# scripts/etc/systemd/system/qmanager-console.service
# /lib/systemd/system/qmanager-console.service
[Unit]
Description=QManager Web Console (ttyd)
After=network.target lighttpd.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usrdata/qmanager/console/ttyd -i 127.0.0.1 -p 8080 -t 'theme={"foreground":"#e4e4e7","background":"#09090b","cursor":"#e4e4e7"}' -t fontSize=14 --writable /usrdata/qmanager/console/console.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
