# scripts/etc/systemd/system/qmanager-firewall.service
# /lib/systemd/system/qmanager-firewall.service
[Unit]
Description=QManager Port Firewall
Before=qmanager-setup.service lighttpd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/qmanager_firewall start
ExecStop=/usr/bin/qmanager_firewall stop

[Install]
WantedBy=multi-user.target
